# systemd/system/heat-api.service
[Unit]
Description=OpenStack Heat API Server
After=syslog.target network.target

[Service]
Type=notify
NotifyAccess=all
TimeoutStartSec=0
Restart=on-failure
User=root
ExecStart=/usr/bin/heat-api --config-file=/etc/heat/heat.conf --log-file=/var/log/heat/heat-api.log
ExecStop=/bin/kill -WINCH ${MAINPID}
KillSignal=SIGCONT

[Install]
WantedBy=multi-user.target

